# x11/xconfig/pkg-plist
bin/xconfig
etc/rc.d/xconfig
etc/X11/cardDetect/XF86Config.amdgpu
etc/X11/cardDetect/XF86Config.bhyve
etc/X11/cardDetect/XF86Config.dual
etc/X11/cardDetect/XF86Config.hyperv
etc/X11/cardDetect/XF86Config.intel
etc/X11/cardDetect/XF86Config.kaveri
etc/X11/cardDetect/XF86Config.qemu
etc/X11/cardDetect/XF86Config.radeonkms
etc/X11/cardDetect/XF86Config.safe
etc/X11/cardDetect/XF86Config.scfb
etc/X11/cardDetect/XF86Config.vesa
etc/X11/cardDetect/XF86Config.virtualbox
etc/X11/cardDetect/XF86Config.vmware
