# x11/nvidia-hybrid-graphics/files/openrc-nvidia_xorg.in
#!/sbin/openrc-run

. /etc/rc.subr

name=nvidia_xorg
desc="Nvidia GPU Dedicated Xorg"

load_rc_config $name

: ${nvidia_modeset:=NO}
: ${nvidia_manage_kld:=YES}

command=/usr/sbin/daemon
pidfile=/var/run/${name}.pid
procname="%%LOCALBASE%%/bin/Xorg"
command_args="-p ${pidfile} -f %%PREFIX%%/bin/Xorg-nvidia-headless"

depend(){
        after localmount devd
        keyword -shutdown
}

if checkyesno nvidia_modeset; then
	kmod=nvidia-modeset
else
	kmod=nvidia
fi

start()
{
	if checkyesno nvidia_manage_kld; then
		kldload -n ${kmod} || return 1
	fi
}

stop()
{
	if checkyesno nvidia_manage_kld; then
		kldunload -f ${kmod} || return 1
	fi
}
