# emulators/open-vm-tools/files/openrc-vmware-guestd.in
#!/sbin/openrc-run

# VMware guest daemon
name="vmware_guestd"
command="/usr/local/bin/vmtoolsd"
command_args="-c /usr/local/share/vmware-tools/tools.conf -p /usr/local/lib/open-vm-tools/plugins/vmsvc"
pidfile="/var/run/${name}.pid"

start() {
        ebegin "Starting $name"
        start-stop-daemon --start --exec ${command} -b -m -p ${pidfile} -- ${command_args}
        eend $?
}

stop() {
        ebegin "Stopping $name"
        start-stop-daemon --stop --exec ${command} -p ${pidfile}
        eend $?
}
